# game-server/script/crontab
0 0 * * * mongo < ~/workspace/dist/zgz/game-server/script/rankingTask.js